# scripts/kb-proxy.service
[Unit]
Description=grimoire REST server
After=network.target

[Service]
Type=simple
# TODO Phase 3: set WorkingDirectory, EnvironmentFile, and ExecStart for your host
WorkingDirectory=/path/to/grimoire
EnvironmentFile=/path/to/grimoire/.env
ExecStart=/path/to/grimoire/.venv/bin/uvicorn server.app:app --host 0.0.0.0 --port 8000
Restart=on-failure

[Install]
WantedBy=multi-user.target
